# archiso/airootfs/etc/systemd/system/calamares.service
[Unit]
Description=Mount directory for calamares configs

[Service]
Type=oneshot
ExecStart=sudo mount -o ro,bind /etc/calamaresmod /etc/calamares
